# wizfi630s_package/wizfi630s-init/files/etc/hotplug.d/button/buttons
#!/bin/sh
logger the hotplug.d/buttons/button was $BUTTON and the action was $ACTION

echo $BUTTON $ACTION > /dev/console 

