# wizfi630s_package/wizfi630s-init/files/etc/init.d/wizfi630s_init
#!/bin/sh /etc/rc.common

START=90

start()
{
	pinmux set i2c gpio
}
