# wizfi630s_package/wizfi630s-init/files/etc/uci-defaults/69_wizfi630s_end
#!/bin/sh

uci set system.@system[0].fb=1
uci commit system
/etc/init.d/network restart
